# builder/guest/etc/systemd/system/startdir.mount
[Unit]
Description=Mount /startdir with virtiofs

[Mount]
What=startdir
Where=/startdir
Type=virtiofs
